# roles/keepalived/templates/keepalived_stats_exporter.service.j2
[Unit]
Description=keepalived_stats_exporter service
After=network-online.target

[Service]
User=root
ExecStart=/usr/local/bin/keepalived-exporter

[Install]
WantedBy=multi-user.target
